# supervisord.conf
[supervisord]

nodaemon=true

[program:apache2]

command=/start2.sh

[program:mysql]

command=/start.sh mysqld
